# home/private_dot_config/systemd/user/mbsync.timer
[Unit]
Description=Mailbox synchronization timer

[Timer]
OnBootSec=2m
OnUnitActiveSec=10m
Unit=mbsync.service

[Install]
WantedBy=timers.target
